# homeworks/hw03/ansible/files/backend.service.j2
[Unit]
Description=My Custom Service
After=network.target

[Service]
WorkingDirectory=/usr/local/be
Environment=ASPNETCORE_URLS=http://*:5001
ExecStart=/usr/local/be/realworld.api
Restart=always